# configs/dot_gateway-0.0.1.conf
[program:{{.Program}}]
user={{.User}}
directory={{.UserHome}}
command={{.UserHome}}/.marlin/ctl/bin/dot_gateway-{{.Version}}
priority=100
autostart=true
autorestart=true
