# configs/near_gateway-0.0.1.conf
[program:{{.Program}}]
user={{.User}}
directory={{.UserHome}}
command={{.UserHome}}/.marlin/ctl/bin/near_gateway-{{.Version}}
priority=100
autostart=true
autorestart=true
